# packages/plusdeck/plusdeck-3.0.0.tar.gz/plusdeck-3.0.0/systemd/plusdeck.service
[Unit]
Description=Plus Deck 2C PC Cassette Deck

[Service]
Type=dbus
BusName=org.jfhbrook.plusdeck
ExecStart=/usr/bin/plusdeckd
Restart=on-failure

